# scripts/kvittra-webhook.service
[Unit]
Description=Kvittra Auto-Deploy Webhook
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/lvcmediahub
Environment=DEPLOY_SECRET=ÄNDRA_MIG
Environment=DEPLOY_BRANCH=claude/explore-migrate-supabase-gGM4t
Environment=DEPLOY_PORT=9090
ExecStart=/usr/bin/node /opt/lvcmediahub/scripts/webhook-server.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
